# python_scripts/py-autoreg.ini
[program:py-autoreg]
command=/bin/bash -lc "/usr/local/bin/auto-register.sh"
directory=/workspace
autostart=true
autorestart=true
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/py-autoreg.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/py-autoreg.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=SUPERVISOR_USER="%(ENV_SUPERVISOR_USER)s",SUPERVISOR_PASSWORD="%(ENV_SUPERVISOR_PASSWORD)s"
